# Maintenance/crontab
0 0 * * * /etc/webmin/package-updates/update.pl
0 0 * * * /home/server/maintenance.sh